# deploy/systemd/gcd-inbound.service
[Unit]
Description=Green Candle Dispatch Inbound Listener
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/srv/gcloads-app
EnvironmentFile=/srv/gcloads-app/.env
ExecStart=/srv/gcloads-app/.venv/bin/python /srv/gcloads-app/inbound_listener.py
Restart=always
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target